# libs-pistrap/raspi-setup.service
[Unit]
Description=Raspi-Setup

[Service]
Type=oneshot
ExecStart=/etc/profile.d/raspi-config.sh

[Install]
WantedBy=multi-user.target